# debian/repowerd.upstart
author "Alexandros Frantzis <alexandros.frantzis@canonical.com>"

start on started dbus and android

respawn

exec /usr/sbin/repowerd
